# supervisord.conf
[supervisord]
nodaemon=true

[program:polipo]
command=su-exec root /usr/bin/polipo -c /ssr-data/polipo/config

[program:ssr]
command=su-exec root python /opt/shadowsocks/local.py -c /ssr-data/shadowsocks/config.json
